# deploy-package/ienet.service
[Unit]
Description=IeNet React Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/ienet.online
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
